# deploy/portal-admin-api.service
[Unit]
Description=Portal ZTNA Admin API
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/portal-ztna-luckygaming
Environment=NODE_ENV=production
Environment=ADMIN_API_PORT=19191
ExecStart=/usr/bin/npm run admin-api:start
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
